# contrib/debian/hibenchmarks.service
[Unit]
Description=hibenchmarks real-time system monitoring
After=network.target

[Service]
Type=simple
EnvironmentFile=-/etc/default/hibenchmarks
ExecStart=/usr/sbin/hibenchmarks -D $EXTRA_OPTS
TimeoutStopSec=30
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
